# provisioning/roles/external.benchmark-worker/files/etc/systemd/system/benchmark-worker.service
[Unit]
Description = isucon9 qualifier benchmark worker

[Service]
WorkingDirectory=/home/isucon/isucari
LimitNOFILE=655360

ExecStart = /home/isucon/isucari/bin/benchmark-worker -ep https://portal.isucon.net

Restart   = always
Type      = simple
User      = isucon
Group     = isucon

[Install]
WantedBy = multi-user.target
